# root/etc/s6-overlay/s6-rc.d/svc-mod-transmission-mod/run
#!/usr/bin/with-contenv bash

# Run magnet watch in background
/usr/local/bin/magnet.sh